# ansible/roles/scoreboard/files/mlb-scoreboard.conf
[program:scoreboard]
command=python3 main.py --led-gpio-mapping=adafruit-hat-pwm --led-brightness=40 --led-slowdown-gpio=3  --led-cols=64
directory=/opt/mlb-led-scoreboard
autostart=true
autorestart=true
