# ansible/roles/scannerpi/templates/node.supervisor.conf.j2
[program:scanner-service]
directory={{ root_scanner_dir }}/node
command=npm run scanner
autostart=true
autorestart=true
user={{ ansible_user }}

