# ansible/roles/scannerpi/templates/node.supervisor.conf.j2
[program:scanner-service]
directory={{ root_scanner_dir }}
command={{ root_scanner_dir }}/scanner-api
autostart=true
autorestart=true
user={{ ansible_user }}
